# packaging/systemd/neonhud.service
[Unit]
Description=NeonHud Performance HUD
Documentation=https://github.com/yourusername/NeonHud
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/neonhud dash
Restart=on-failure
RestartSec=5
User=%i
WorkingDirectory=/home/%i

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
